# devices/merry/root/etc/systemd/system/node_exporter.service
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=prometheus
Group=root
Type=simple
ExecStart=/opt/node_exporter/live/node_exporter

[Install]
WantedBy=multi-user.target
